# telegram-proxy.service
[Unit]
Description=Iron Lady Assistant Telegram Read Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-developer
WorkingDirectory=/home/claude-developer/iron-lady-assistant
ExecStart=/home/claude-developer/iron-lady-assistant/venv/bin/python3 -m src.telegram_proxy
Restart=always
RestartSec=5
EnvironmentFile=/home/claude-developer/iron-lady-assistant/.env
Environment=PATH=/home/claude-developer/.npm-claude-developer/bin:/home/claude-developer/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Prefer a systemd-encrypted credential blob here. Example:
#   sudo systemd-creds encrypt --name=telegram_proxy_key - /etc/credstore.encrypted/telegram_proxy_key < /root/telegram_proxy_key
# Then keep this line aligned with TELEGRAM_PROXY_KEY_CREDENTIAL_NAME in .env.
LoadCredentialEncrypted=telegram_proxy_key:/etc/credstore.encrypted/telegram_proxy_key

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/claude-developer/iron-lady-assistant/memory
PrivateTmp=true

[Install]
WantedBy=multi-user.target
