# telegram-bot.service
[Unit]
Description=Telegram Claude Code Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-developer
WorkingDirectory=/home/claude-developer/claude-code-as-assistant
ExecStart=/home/claude-developer/claude-code-as-assistant/run.sh
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=600
EnvironmentFile=/home/claude-developer/claude-code-as-assistant/.env
Environment=PATH=/home/claude-developer/.npm-claude-developer/bin:/home/claude-developer/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/claude-developer
PrivateTmp=true

[Install]
WantedBy=multi-user.target
